# example.automount
[Unit]
Description=Automount example

[Automount]
Where=/mnt/example
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
